# packer/linux/configs/firstboot.service
[Unit]
Description=First boot setup script

[Service]

[Install]